# systemd/.config/systemd/user/task-sync.timer
[Unit]
Description=Regularly sync tasks for TaskWarrior

[Timer]
OnStartupSec=30min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
